# deploy/trip-report-backup.service
[Unit]
Description=Trip Report Database Backup

[Service]
Type=oneshot
User=postgres
ExecStart=/opt/trip-report/deploy/backup.sh
